# scripts/football-platform.service
[Unit]
Description=Football Prediction Platform — Docker Compose Services
Requires=network-online.target
After=network-online.target docker.service
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=oneshot
RemainAfterExit=yes
User=diego
WorkingDirectory=/home/diego/EngenhariaDeDados

ExecStart=/bin/bash /home/diego/EngenhariaDeDados/scripts/start_all.sh
ExecStop=/bin/bash /home/diego/EngenhariaDeDados/scripts/stop_all.sh

TimeoutStartSec=300
TimeoutStopSec=60

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
